# cspp-update.cron
#
# Initial run only???
#sdr_luts.sh -d -W /home/runuser/tmp
#
# Daily update required for CSPP v3.1.2
SHELL=/bin/bash
0 3 * * * runuser source /home/runuser/.bashrc; sdr_ancillary.sh -d -W /home/runuser/tmp
